# bin/com.jonmagic.homeassistant-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.jonmagic.homeassistant-tunnel</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/jonmagic/.dotfiles/bin/tsp-homeassistant</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>30</integer>
	<key>StandardOutPath</key>
	<string>/tmp/homeassistant-tunnel.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/homeassistant-tunnel.log</string>
</dict>
</plist>
